# debian/systemd/koha.target
[Unit]
Description=Koha ILS - All instances
After=network.target memcached.service mysql.service

[Install]
WantedBy=multi-user.target
