# debian/koha-sysv.prerm
#!/bin/sh

set -e

if [ "$1" = "remove" ]; then
    # Stop the service
    if [ -d /run/systemd/system ]; then
        deb-systemd-invoke stop koha-common.service >/dev/null || true
    else
        invoke-rc.d koha-common stop || true
    fi
fi

#DEBHELPER#

exit 0
